# CONFIGURATION/etc/rc.local
while true; do bash /root/hangup_chan.sh; sleep 20; done
